# scripts/com.luckyhands.intel_daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.luckyhands.intel_daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ritzmish/intel_daily/scripts/run_daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ritzmish/intel_daily</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ritzmish/intel_daily/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ritzmish/intel_daily/logs/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
